# bin/bot.service
[Unit]
Description=Telegram bot to resolve user expressions
After=default.target

[Service]
User=root
Restart=on-failure
ExecStart=/usr/bin/tg-autoresolver

[Install]
WantedBy=network.target